# launch-control/set-meeting-light.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>set-meeting-light</string>
	<key>ProcessType</key>
	<string>Background</string>
	<key>Program</key>
	<string>/Users/jcsims/bin/set-meeting-light</string>
	<key>StandardErrorPath</key>
	<string>/tmp/set-meeting-light-err.log</string>
	<key>StandardOutPath</key>
	<string>/tmp/set-meeting-light-out.log</string>
	<key>StartCalendarInterval</key>
	<array>
		<dict/>
	</array>
</dict>
</plist>
